# portable-walkthrough-go.socket
[Unit]
Description=Portable Walkthrough Go Edition Socket

[Socket]
ListenStream=8080

[Install]
WantedBy=sockets.target
